# system_files/usr/lib/systemd/user/gamescope-session-plus@steam.service.d/10-armada-powerbuttond.conf
[Unit]
Wants=armada-powerbuttond.service
After=armada-powerbuttond.service
