# system_files/usr/lib/systemd/system/ostree-finalize-staged.service.d/10-armada-bootimg.conf
[Unit]
RequiresMountsFor=/boot/efi

[Service]
# Appended after OSTree's own ExecStop=, so /KERNEL is generated only after a
# staged deployment has finalized successfully.
ExecStop=/usr/libexec/armada/armada-bootimg-finalize
